# librephotos/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
pidfile=/var/run/supervisord.pid
user=root

[program:postgresql]
command=gosu postgres /usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/data -c config_file=/etc/postgresql/15/main/postgresql.conf
autostart=true
autorestart=true
stdout_logfile=/var/log/postgresql.log
stderr_logfile=/var/log/postgresql_err.log
priority=100

[program:redis]
command=/usr/bin/redis-server --daemonize no --bind 127.0.0.1
autostart=true
autorestart=true
stdout_logfile=/var/log/redis.log
stderr_logfile=/var/log/redis_err.log
priority=110

[program:librephotos-backend]
command=/code/entrypoint.sh
directory=/code
environment=
    DB_BACKEND=postgresql,
    DB_NAME=librephotos,
    DB_USER=librephotos,
    DB_PASS=%(ENV_DB_PASS)s,
    DB_HOST=127.0.0.1,
    DB_PORT=5432,
    REDIS_HOST=127.0.0.1,
    REDIS_PORT=6379,
    BACKEND_HOST=127.0.0.1,
    BACKEND_PORT=8001,
    SECRET_KEY=%(ENV_SECRET_KEY)s,
    ALLOWED_HOSTS=*,
    DATA_ROOT=%(ENV_SCAN_DIRECTORY)s
autostart=true
autorestart=true
stdout_logfile=/var/log/librephotos.log
stderr_logfile=/var/log/librephotos_err.log
priority=200

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx_access.log
stderr_logfile=/var/log/nginx_error.log
priority=300
