# librephotos/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
pidfile=/var/run/supervisord.pid
user=root

[program:postgresql]
command=/usr/bin/postgres -D /var/lib/postgresql/data
user=postgres
autostart=true
autorestart=true
stdout_logfile=/var/log/postgresql.log
stderr_logfile=/var/log/postgresql_err.log
priority=100

[program:redis]
command=/usr/bin/redis-server --daemonize no
autostart=true
autorestart=true
stdout_logfile=/var/log/redis.log
stderr_logfile=/var/log/redis_err.log
priority=110

[program:librephotos-worker]
command=/opt/librephotos/venv/bin/celery -A librephotos worker -l %(ENV_DJANGO_LOG_LEVEL)s
directory=/opt/librephotos/repo/apps/backend
user=librephotos
autostart=true
autorestart=true
stdout_logfile=/var/log/celery.log
stderr_logfile=/var/log/celery_err.log
priority=200

[program:librephotos-backend]
command=/opt/librephotos/venv/bin/gunicorn \
    --bind 127.0.0.1:8001 \
    --workers %(ENV_WORKERS)s \
    --timeout 600 \
    --error-logfile /var/log/gunicorn_error.log \
    librephotos.wsgi:application
directory=/opt/librephotos/repo/apps/backend
user=librephotos
autostart=true
autorestart=true
stdout_logfile=/var/log/gunicorn.log
stderr_logfile=/var/log/gunicorn_err.log
priority=200

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx/access.log
stderr_logfile=/var/log/nginx/error.log
priority=300
